# support/mqtt/wardrive-maint.timer
[Unit]
Description=MeshCore Coverage Maintenance

[Timer]
OnCalendar=Hourly
Persistent=true

[Install]
WantedBy=timers.target